# oanda-bot.service
[Unit]
Description=OANDA Telegram Signals Bot
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/opt/oanda-telegram-bot
Environment="PATH=/opt/oanda-telegram-bot/venv/bin"
ExecStart=/opt/oanda-telegram-bot/venv/bin/python -m src.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

